# crontabs/raspi2.latest
# 12 * * * * /home/pi/capture/bin/upload.sh >> /home/pi/capture/logs/upload.log 2>&1
5,10,15,20,25,30,35,40,45,50,55 * * * *  /home/pi/capture/bin/watchdog.sh >> /home/pi/capture/logs/watchdog.log 2>&1
